# printerbot.service
[Unit]
Description=Telegram Printer and Scanner Bot
After=network.target

[Service]
User=username
WorkingDirectory=/home/username
ExecStart=/usr/bin/python3 /home/username/businessbot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
